# luci-app-nfs-server/root/etc/init.d/nfs-server
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=0

generate_exports() {
    echo "# Autogenerated by LuCI NFS Server" > /etc/exports
    for section in $(uci show nfs | grep "=share" | cut -d. -f2 | cut -d= -f1); do
        path=$(uci get nfs.$section.path 2>/dev/null)
        opts=$(uci get nfs.$section.options 2>/dev/null)
        [ -n "$path" ] || continue
        [ -n "$opts" ] || opts="rw,sync,no_subtree_check"
        echo "$path *( $opts )" >> /etc/exports
    done
}

start() {
    generate_exports
    /etc/init.d/nfsd restart 2>/dev/null || /etc/init.d/nfs-kernel-server restart 2>/dev/null
}
